# lcdpi/modconf/ztermcom.socket
Unit]
Description=Z-Terminal Communication Socket

[Socket]
ListenFIFO=/dev/zterm
Service=ztermcom.service
