# lcdpi/modconf/ztermcom.socket
Unit]
Description=Z-Terminal Communication Socket

[Socket]
Backlog=4
MaxConnections=1
ListenFIFO=/tmp/zterm
Service=ztermcom.service
