# rpi/config/systemd/rpi-cleanup.timer
[Unit]
Description=Runs '/opt/rpi/cleanup' script every 12 hours

[Timer]
OnBootSec=2h
OnActiveSec=12h

[Install]
WantedBy=timers.target