# rpi/config/systemd/rpi-netdetect.timer
[Unit]
Description=Runs script every 5min.

[Timer]
OnBootSec=30
OnUnitActiveSec=5min

[Install]
WantedBy=timers.target