# rpi/config/systemd/rpi-fcleds.service
[Unit]
Description=Fadecandy LED Effects Loop

[Service]
User=root
Type=notify
ExecStart=/opt/rpi/leds.sysd
Restart=on-failure
##KillSignal=SIGKILL
RestartSec=5
