# rpi/config/systemd/rpi-rebootpi.timer
[Unit]
Description=Restart Pi six days

[Timer]
OnActiveSec=144h

[Install]
WantedBy=timers.target