# rpi/config/systemd/rpi-rebootpi.timer
[Unit]
Description=Restart Pi every six days

[Timer]
OnActiveSec=6d

[Install]
WantedBy=timers.target