# rpi/config/systemd/rpi-camcleanup.service
[Unit]
Description=Camera Cleanup

[Service]
Type=oneshot
RemainAfterExit=no
ExecStart=-/opt/rpi/camcleanup
