# rpi/config/systemd/rpi-cleanup.service
[Unit]
Description=12hr Clean-up Script

[Service]
Type=oneshot
RemainAfterExit=no
ExecStart=-/opt/rpi/cleanup
IOSchedulingClass=idle
