# rpi/config/systemd/rpi-fcserver.service
[Unit]
Description=Fadecandy Server
Conflicts=shutdown.target

[Service]
User=root
ExecStart=/bin/bash -c "/opt/rpi/effects/fcserver /opt/rpi/effects/$rpi_ledtype/fcbin/$rpi_fcsetup"
WorkingDirectory=/opt/rpi/effects
Restart=always
RestartSec=5
