# rpi/config/systemd/rpi-rebootpi.service
[Unit]
Description=Restart Pi

[Service]
LogLevelMax=notice
Type=oneshot
RemainAfterExit=no
ExecStart=/bin/systemctl reboot
