# campi/modconf/rpi-camcleanup.service
[Unit]
Description=Camera Cleanup

[Service]
Type=oneshot
RemainAfterExit=no
ExecStart=/opt/rpi/camcleanup > /dev/null 2>&1 &
