# homepi/modconf/rpi-motionping.timer
[Unit]
Description=Restart motion if cameras disconnect

[Timer]
OnActiveSec=10min

[Install]
WantedBy=timers.target