# files/config/file-monitor.service
[Unit]
Description=File Monitor API Trigger
After=local-fs.target network.target

[Service]
User=monitor
Group=monitor
ExecStart=/usr/bin/file_monitor
Type=simple
Restart=on-failure
RestartSec=5
