# pve/config/actiontrig.timer
[Unit]
Description=Triggers script every 40 seconds.

[Timer]
OnBootSec=5min
OnUnitActiveSec=40

[Install]
WantedBy=timers.target